# 1password/.config/autostart/1password.desktop
[Desktop Entry]
Type=Application
Name=1Password
Exec=/usr/bin/1password --silent
Hidden=false
X-GNOME-Autostart-enabled=true
